# services/base/sudo-invalidate-on-lock.service
[Unit]
Description=Invalidate sudo credentials on screen lock

[Service]
ExecStart=%h/.local/bin/sudo-invalidate-on-lock
Restart=on-failure

[Install]
WantedBy=default.target
